# unit_files/uiowa_gym.service
[Unit]
Description=Automatically register for a gym slot

[Service]
Type=oneshot
ExecStart=/usr/bin/uiowa-gym
